# etc-power-astute/etc/systemd/system/astute-idle-suspend.service
[Unit]
Description=Astute idle suspend check

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/astute-idle-check.sh

