# backup-systemd-audacious/etc/systemd/system/borg-check-deep.service
[Unit]
Description=Monthly deep Borg data integrity check (audacious -> astute)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=alchemist
Group=alchemist
WorkingDirectory=/home/alchemist
SuccessExitStatus=0 1

EnvironmentFile=/home/alchemist/.config/borg/env

# wake NAS + wait
ExecStartPre=/usr/local/lib/borg/wait-for-astute.sh

ExecStartPre=/usr/bin/test -f /home/alchemist/.config/borg/passphrase

# Full verification of repository data blocks (can be slow on large repos)
ExecStart=/usr/bin/borg check --verbose --verify-data "$BORG_REPO"

# Give it plenty of time, since this is allowed to be heavyweight
TimeoutStartSec=14400
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
