# backup-systemd-audacious/etc/systemd/system/borg-check-deep.service
[Unit]
Description=Monthly deep Borg data integrity check (audacious -> astute)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=alchemist
Group=alchemist
WorkingDirectory=/home/alchemist
SuccessExitStatus=0 1

# Repository connection + credentials (same as borg-backup.service / borg-check.service)
Environment="BORG_REPO=borg@192.168.1.154:/srv/backups/audacious-borg"
Environment="BORG_RSH=ssh -i /home/alchemist/.ssh/audacious-backup -o ConnectTimeout=20 -o ServerAliveInterval=30 -o ServerAliveCountMax=120"
Environment="BORG_PASSCOMMAND=cat /home/alchemist/.config/borg/passphrase"

# Wake Astute and confirm it’s actually ready
ExecStartPre=/usr/bin/wakeonlan 60:45:cb:9b:ab:3b
ExecStartPre=/usr/bin/ssh -i /home/alchemist/.ssh/audacious-backup borg@192.168.1.154 /usr/bin/mountpoint -q /srv/backups
ExecStartPre=/usr/bin/test -f /home/alchemist/.config/borg/passphrase

# Full verification of repository data blocks (can be slow on large repos)
ExecStart=/usr/bin/borg check --verbose --verify-data

# Give it plenty of time, since this is allowed to be heavyweight
TimeoutStartSec=14400
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target