# cold-storage-audacious/.config/systemd/user/cold-storage-reminder.service
[Unit]
Description=Cold storage backup reminder

[Service]
Type=oneshot
ExecStart=%h/.local/bin/cold-storage-reminder.sh
