# root-power-astute/etc/systemd/system/astute-idle-suspend.service
[Unit]
Description=Astute idle suspend check
After=srv-nas.mount zfs-mount.service
RequiresMountsFor=/srv/nas

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/astute-idle-check.sh

