# root-power-astute/etc/systemd/system/astute-idle-suspend.timer
[Unit]
Description=Run Astute idle suspend check every 3 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=3min
AccuracySec=30s

[Install]
WantedBy=timers.target

