# root-borg-audacious/etc/systemd/system/borg-check-deep.service
[Unit]
Description=Monthly deep Borg data integrity check (audacious -> astute)
Wants=network-online.target
After=network-online.target home.mount zfs-mount.service
RequiresMountsFor=/home/alchemist

[Service]
Type=oneshot
User=alchemist
Group=alchemist
WorkingDirectory=/home/alchemist

Environment="HOME=/home/alchemist"
EnvironmentFile=/home/alchemist/.config/borg/env

# Deep checks must not be killed mid-run
TimeoutStartSec=0
Restart=no

# Wake and wait for astute
ExecStartPre=/usr/local/lib/borg/wait-for-astute.sh

# Verify passphrase file exists
ExecStartPre=/usr/bin/test -f /home/alchemist/.config/borg/passphrase

# Protected from idle shutdown by idle-shutdown.sh checking service status
ExecStart=/usr/local/lib/borg/run-deep-check.sh

Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
