# root-borg-audacious/etc/systemd/system/borg-check-deep.timer
[Unit]
Description=Timer: Monthly deep Borg integrity check (audacious -> astute)

[Timer]
# First Sunday of each month at 11:00 local time
OnCalendar=Sun *-*-01..07 11:00
Persistent=true
RandomizedDelaySec=600

# Retry every 6 hours if the service is inactive (e.g. Astute asleep)
OnUnitInactiveSec=6h

[Install]
WantedBy=timers.target
