# root-offsite-astute/etc/systemd/system/borg-offsite-astute-critical.service
[Unit]
Description=Borg off-site: Astute critical data to BorgBase (append-only access)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/srv/nas /var/lib/borg-offsite /var/cache/borg-offsite

[Service]
Type=oneshot
User=root
Group=root
StandardOutput=journal
StandardError=journal

Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

TimeoutStartSec=0

ExecStart=/usr/bin/systemd-inhibit --what=shutdown:sleep --mode=block /usr/local/lib/borg-offsite/run-astute-critical.sh
