# root-offsite-astute/etc/systemd/system/borg-offsite-astute-critical.timer
[Unit]
Description=Daily Borg off-site push (astute-critical)

[Timer]
OnCalendar=*-*-* 15:00:00
Persistent=true
WakeSystem=true
AccuracySec=5min
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
