# root-backup-audacious/etc/systemd/system/borg-backup.service
[Unit]
Description=BorgBackup: Audacious → Astute NAS (6-hourly, WoL Astute)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/alchemist

Environment="HOME=/home/alchemist"
EnvironmentFile=/home/alchemist/.config/borg/env

SuccessExitStatus=0 1

# Be nice to foreground work
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Fail fast if Astute is asleep; timer will run again next hour
Restart=no
TimeoutStartSec=0

ExecStart=/usr/bin/systemd-inhibit \
  --what=shutdown:sleep \
  --why="Borg backup in progress (Audacious → Astute)" \
  /usr/local/lib/borg/run-backup.sh
