# root-backup-audacious/etc/systemd/system/borg-check-deep.service
[Unit]
Description=Monthly deep Borg data integrity check (audacious -> astute)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/alchemist

Environment="HOME=/home/alchemist"
EnvironmentFile=/home/alchemist/.config/borg/env

# Deep checks must not be killed mid-run
TimeoutStartSec=0
Restart=no

ExecStart=/usr/bin/systemd-inhibit \
  --what=shutdown:sleep \
  --why="Borg deep integrity check in progress" \
  /usr/local/lib/borg/run-deep-check.sh

Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
