# root-backup-audacious/etc/systemd/system/borg-check.service
[Unit]
Description=BorgBackup: Integrity check (Audacious → Astute NAS)
Wants=network-online.target
After=network-online.target home.mount zfs-mount.service
RequiresMountsFor=/home/alchemist
# Run only when Astute is awake / reachable.

[Service]
Type=oneshot
User=alchemist
Group=alchemist
WorkingDirectory=/home/alchemist
SuccessExitStatus=0 1

EnvironmentFile=/home/alchemist/.config/borg/env

# wake and wait
ExecStartPre=/usr/local/lib/borg/wait-for-astute.sh

# local passphrase check
ExecStartPre=/usr/bin/test -f /home/alchemist/.config/borg/passphrase

# The actual integrity check
ExecStart=/usr/bin/borg check --verbose "$BORG_REPO"

# Conservative timeout — large repos can take a while
TimeoutStartSec=3600
