# root-backup-audacious/etc/systemd/system/borg-backup.service
[Unit]
Description=BorgBackup: Audacious → Astute NAS (6-hourly, WoL Astute)
Wants=network-online.target
After=network-online.target home.mount var.mount zfs-mount.service
RequiresMountsFor=/home/alchemist /var/lib/borg /var/cache/borg

[Service]
Type=oneshot
User=alchemist
Group=alchemist
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/alchemist

# System backup uses separate state directories (no conflict with user borg commands)
Environment="BORG_BASE_DIR=/var/lib/borg/audacious"
Environment="BORG_CONFIG_DIR=/var/lib/borg/audacious/config"
Environment="BORG_CACHE_DIR=/var/cache/borg/audacious"
Environment="BORG_SECURITY_DIR=/var/lib/borg/audacious/security"

# Backup script will source passphrase and SSH config from user directory
Environment="HOME=/home/alchemist"

SuccessExitStatus=0 1

# Be nice to foreground work
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Fail fast if Astute is asleep; timer will run again next hour
Restart=no
TimeoutStartSec=0

ExecStart=/usr/local/lib/borg/run-backup-with-logging.sh
