# host-audacious/etc-power/etc/udev/rules.d/69-hdparm.rules
ACTION=="add|change", KERNEL=="sda", RUN+="/usr/sbin/hdparm -S 120 /dev/sda"
