# etc-systemd/system/borg-backup.service
[Unit]
Description=BorgBackup: Audacious → Astute NAS (hourly, deduped)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot

# Repo, SSH key, and passphrase provider
Environment="BORG_REPO=borg@192.168.1.154:/srv/nas/backups/audacious-borg"
Environment="BORG_RSH=ssh -i /home/alchemist/.ssh/audacious-backup -o ConnectTimeout=20"
Environment="BORG_PASSCOMMAND=cat /home/alchemist/.config/borg/passphrase"

# Be nice to foreground work
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Wake Astute before backup (your MAC)
ExecStartPre=/usr/bin/wakeonlan 60:45:cb:9b:ab:3b

# Create snapshot (non-hidden top-level in ~; exclusions via patterns file)
ExecStart=/usr/bin/borg create --verbose --stats --compression lz4 --one-file-system \
  ::{now:%Y-%m-%dT%H:%M:%S} \
  --patterns-from /home/alchemist/.config/borg/patterns

# Prune to your policy: hourly (within 3d), daily(14), weekly(52 ~= 1yr)
ExecStartPost=/usr/bin/borg prune --list \
  --keep-within 3d \
  --keep-daily 14 \
  --keep-weekly 52

# Fail fast if Astute is asleep; timer will run again next hour
Restart=no
TimeoutStartSec=600