# etc-systemd/system/efi-sync.path
[Unit]
Description=Trigger ESP sync when kernel-install updates UKIs

[Path]
PathChanged=/boot/efi/EFI/Linux/
Unit=efi-sync.service

[Install]
WantedBy=multi-user.target

