# nas-audacious/.config/systemd/user/astute-nas.service
[Unit]
Description=Astute NAS session (wake, inhibit, mount)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=%h/.config/systemd/user/astute-nas.env

# Wake NAS provider (only if not already reachable)
ExecStart=/usr/bin/bash -c '\
  if ping -c1 -W1 ${NAS_PROVIDER_HOST} >/dev/null 2>&1; then \
    logger -t astute-nas "${NAS_PROVIDER_HOST} already reachable, skipping WOL"; \
    exit 0; \
  fi; \
  logger -t astute-nas "Sending WOL to ${NAS_PROVIDER_HOST}"; \
  /usr/bin/wakeonlan ${NAS_PROVIDER_MAC}; \
  for i in {1..30}; do \
    ping -c1 -W1 ${NAS_PROVIDER_HOST} >/dev/null 2>&1 && exit 0; \
    sleep 1; \
  done; \
  exit 1'

# Enable sleep inhibit on NAS provider
ExecStart=/usr/bin/ssh astute-nas start

# Mount NAS locally
ExecStart=/usr/bin/sudo /usr/bin/systemctl start srv-astute.mount

# Guaranteed teardown
ExecStop=/usr/bin/sudo /usr/bin/systemctl stop srv-astute.mount
ExecStop=/usr/bin/ssh astute-nas stop

[Install]
WantedBy=default.target
